# supervisord.conf
[supervisord]
nodaemon=true

[program:ziproxy]
command=/usr/bin/ziproxy -d

